# install_files/echoarp.timer
[Unit]
Description=OpenSSH server daemon

[Timer]
Persistent=true 
OnCalendar=*-*-* *:*:00
[Install]
WantedBy=timers.target
